# src/pifan-daemon/util/pifand.service.in
[Unit]
Description="Raspberry Pi CPU fan service."

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
ExecStart=@CMAKE_INSTALL_BINDIR@/pifand
TimeoutSec=20
